# packages/sandbox/supervisord.conf
; ==============================================================================
; Airlock Sandbox Supervisor Configuration
; ==============================================================================
; Manages the lifecycle of sandbox processes:
;   1. Xvfb (headless X11 display)
;   2. KasmVNC (WebSocket VNC server)
;   3. Openbox (window manager)
;   4. Target application (Chromium/Evince)
;
; Note: In air-gapped mode (NetworkMode: none), network-dependent features
; are disabled by the container runtime, not supervisord.
; ==============================================================================

[supervisord]
nodaemon=true
user=airlock
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid
loglevel=info

; Reduce log noise in production
childlogdir=/tmp
logfile_maxbytes=1MB
logfile_backups=2

; ==============================================================================
; X11 Virtual Framebuffer
; ==============================================================================
[program:xvfb]
command=Xvfb :1 -screen 0 1920x1080x24 -ac +extension GLX +render -noreset -novtswitch
autorestart=true
priority=10
startsecs=2
startretries=3
stdout_logfile=/tmp/Xvfb.log
stderr_logfile=/tmp/Xvfb.log
environment=DISPLAY=:1

; ==============================================================================
; KasmVNC Server
; ==============================================================================
[program:kasmvnc]
command=vncserver :1 -interface 0.0.0.0 -WebsocketPort 6901 -geometry 1920x1080x24 -MaxFrameRate 24 -PreferBandwidth -DynamicQualityMin=4 -DynamicQualityMax=7 -disableBasicAuth
autorestart=true
priority=20
startsecs=3
startretries=5
stdout_logfile=/tmp/KasmVNC.log
stderr_logfile=/tmp/KasmVNC.log
environment=DISPLAY=:1,HOME=/home/airlock,USER=airlock
; Wait for Xvfb
dependent_startup=true

; ==============================================================================
; Openbox Window Manager
; ==============================================================================
[program:openbox]
command=openbox
autorestart=true
priority=30
startsecs=1
stdout_logfile=/tmp/openbox.log
stderr_logfile=/tmp/openbox.log
environment=DISPLAY=:1
; Wait for KasmVNC
dependent_startup=true

; ==============================================================================
; Target Application Launcher
; ==============================================================================
; This is a one-shot script that launches the appropriate app based on
; TARGET_FILE or TARGET_URL environment variables.
[program:target-launcher]
command=/usr/local/bin/launch-target.sh
autorestart=false
priority=40
startsecs=3
stdout_logfile=/tmp/target.log
stderr_logfile=/tmp/target.log
; Inherit container env (TARGET_FILE, TARGET_URL, DISPLAY, etc.)
; Wait for window manager
dependent_startup=true

; ==============================================================================
; Event Listeners
; ==============================================================================
[eventlistener:process-cleanup]
command=/usr/local/bin/cleanup.sh
events=PROCESS_STATE_STOPPED,PROCESS_STATE_EXITED,PROCESS_STATE_FATAL
buffer_size=10
priority=100

; ==============================================================================
; Group Definition (for collective control)
; ==============================================================================
[group:airlock]
programs=xvfb,kasmvnc,openbox,target-launcher
priority=10
