# packages/sandbox/supervisord.conf
; ==============================================================================
; Airlock Sandbox Supervisor Configuration
; ==============================================================================
; KasmVNC vncserver owns display :1 (X11 + WebSocket VNC). Openbox and the
; target launcher start after the VNC port is listening.
; ==============================================================================

[supervisord]
nodaemon=true
user=airlock
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid
loglevel=info
childlogdir=/tmp
logfile_maxbytes=1MB
logfile_backups=2

[program:kasmvnc]
command=/usr/local/bin/start-kasmvnc.sh
autorestart=true
priority=10
startsecs=5
startretries=5
stdout_logfile=/tmp/KasmVNC.log
stderr_logfile=/tmp/KasmVNC.log
environment=DISPLAY=:1,HOME=/home/airlock,USER=airlock,VNC_PORT=6901,VNC_RESOLUTION=1920x1080,VNC_COL_DEPTH=24,MAX_FRAME_RATE=24

[program:openbox]
command=/bin/sh -c "until DISPLAY=:1 xdpyinfo >/dev/null 2>&1; do sleep 0.5; done; exec openbox"
autorestart=true
priority=20
startsecs=1
stdout_logfile=/tmp/openbox.log
stderr_logfile=/tmp/openbox.log
environment=DISPLAY=:1

[program:target-launcher]
command=/bin/sh -c "until nc -z localhost 6901; do sleep 0.5; done; exec /usr/local/bin/launch-target.sh"
autorestart=false
priority=40
startsecs=3
stdout_logfile=/tmp/target.log
stderr_logfile=/tmp/target.log

[eventlistener:process-cleanup]
command=/usr/local/bin/cleanup.sh
events=PROCESS_STATE_STOPPED,PROCESS_STATE_EXITED,PROCESS_STATE_FATAL
buffer_size=10
priority=100

[group:airlock]
programs=kasmvnc,openbox,target-launcher
priority=10
